# kuso-subdomain-adder-update.timer
[Unit]
Description=Auto Update kuso-subdomain-adder

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
